# config/includes.chroot_after_packages/etc/profile.d/99-kingdom-wallpaper.sh
#!/bin/bash
# Forces the Plasma desktop to use the Kingdom Throne wallpaper

dbus-send --session --dest=org.kde.plasmashell --type=method_call /PlasmaShell org.kde.PlasmaShell.evaluateScript 'string:
var allDesktops = desktops();
for (i=0;i<allDesktops.length;i++) {
    d = allDesktops[i];
    d.wallpaperPlugin = "org.kde.image";
    d.currentConfigGroup = Array("Wallpaper", "org.kde.image", "General");
    d.writeConfig("Image", "file:///usr/share/alfred-linux/kingdom-media/wallpapers/4k/kingdom-throne-4k.png");
}'