# config/includes.chroot_after_packages/etc/systemd/system/qgsm-node.service
[Unit]
Description=QGSM Post-Quantum Network Daemon
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/qgsm-node
Restart=always
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
