# config/includes.chroot_after_packages/etc/skel/.config/autostart/aura-hud.desktop
[Desktop Entry]
Type=Application
Exec=chromium --kiosk --app=file:///opt/aura-dashboard/index.html --window-position=right --window-size=400,1080 --class="AuraHUD"
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=Aura Dashboard HUD
Comment=Unkillable mesh monitoring widget
